# contrib/systemd/bililive-go@.service
[Unit]
Description=Bililive-go record for user %i
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/var/log/
ExecStart=/usr/bin/bililive-go -c /etc/bililive-go/config.yml
Restart=on-failure

[Install]
WantedBy=multi-user.target
